# debian/postinst
#!/bin/sh

set -e

# FIXME: remove in favor of dh_installsystemd
systemctl unmask opensnitch.service
systemctl enable opensnitch.service
service opensnitch restart
