# debian/prerm
#!/bin/sh

set -e

service opensnitch stop || true
systemctl disable opensnitch.service || true
